# rstudio-server.conf
[program:rstudio-server]
command=/bin/bash -c "/usr/lib/rstudio-server/bin/rserver && tail -F /var/log/httpd/access_log"
